# config/ict_trading_system.service.example
# Systemd service unit for ICT Trading System (edit paths as needed)
[Unit]
Description=ICT Smart Money Trading Alert System
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/home/franklin/SOFTWARE_ENGENEERING/Development/code/se-prep/Webportfolio/MYAI-AGENT/prototypeictalgo
ExecStart=/usr/local/bin/uvicorn ict_trading_system.main:app --host 0.0.0.0 --port 8000
EnvironmentFile=/home/franklin/SOFTWARE_ENGENEERING/Development/code/se-prep/Webportfolio/MYAI-AGENT/prototypeictalgo/.env
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
